# Windows/Forensicator-Share/rules/windows/sigma/process_creation/proc_creation_win_schtasks_schedule_via_masqueraded_xml_file.json
{
  "title": "Suspicious Scheduled Task Creation via Masqueraded XML File",
  "rule_id": "dd2a821e-3b07-4d3b-a9ac-929fe4c6ca0c",
  "description": "Detects the creation of a scheduled task using the \"-XML\" flag with a file without the '.xml' extension. This behavior could be indicative of potential defense evasion attempt during persistence",
  "level": "medium",
  "status": "test",
  "tags": [
    "attack.privilege-escalation",
    "attack.execution",
    "attack.defense-evasion",
    "attack.persistence",
    "attack.t1036.005",
    "attack.t1053.005"
  ],
  "category": "process_creation",
  "rule_file": "rules/windows/process_creation/proc_creation_win_schtasks_schedule_via_masqueraded_xml_file.yml",
  "sources": [
    "sysmon_process_creation"
  ],
  "items": {
    "selection_img": {
      "type": "any",
      "children": [
        {
          "type": "all",
          "children": [
            {
              "type": "field",
              "field": "Image",
              "operator": "endswith",
              "values": [
                "\\schtasks.exe"
              ],
              "match": "any",
              "windash": false,
              "ignore_case": true
            }
          ]
        },
        {
          "type": "all",
          "children": [
            {
              "type": "field",
              "field": "OriginalFileName",
              "operator": "exact",
              "values": [
                "schtasks.exe"
              ],
              "match": "any",
              "windash": false,
              "ignore_case": true
            }
          ]
        }
      ]
    },
    "selection_cli_create": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/create",
            "-create"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "selection_cli_xml": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/xml",
            "-xml"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "filter_main_extension_xml": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            ".xml"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "filter_main_system_process": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "IntegrityLevel",
          "operator": "exact",
          "values": [
            "System",
            "S-1-16-16384"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "filter_main_rundll32": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "ParentImage",
          "operator": "endswith",
          "values": [
            "\\rundll32.exe"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "ParentCommandLine",
          "operator": "contains",
          "values": [
            ":\\WINDOWS\\Installer\\MSI",
            ".tmp,zzzzInvokeManagedCustomActionOutOfProc"
          ],
          "match": "all",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "filter_optional_third_party": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "ParentImage",
          "operator": "endswith",
          "values": [
            ":\\ProgramData\\OEM\\UpgradeTool\\CareCenter_*\\BUnzip\\Setup_msi.exe",
            ":\\Program Files\\Axis Communications\\AXIS Camera Station\\SetupActions.exe",
            ":\\Program Files\\Axis Communications\\AXIS Device Manager\\AdmSetupActions.exe",
            ":\\Program Files (x86)\\Zemana\\AntiMalware\\AntiMalware.exe",
            ":\\Program Files\\Dell\\SupportAssist\\pcdrcui.exe"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    }
  },
  "condition": {
    "type": "all",
    "children": [
      {
        "type": "wildcard_ref",
        "mode": "all",
        "pattern": "selection_*"
      },
      {
        "type": "not",
        "child": {
          "type": "wildcard_ref",
          "mode": "any",
          "pattern": "filter_main_*"
        }
      },
      {
        "type": "not",
        "child": {
          "type": "wildcard_ref",
          "mode": "any",
          "pattern": "filter_optional_*"
        }
      }
    ]
  }
}