# Windows/Forensicator-Share/rules/windows/sigma/process_creation/proc_creation_win_schtasks_appdata_local_system.json
{
  "title": "Suspicious Schtasks Execution AppData Folder",
  "rule_id": "c5c00f49-b3f9-45a6-997e-cfdecc6e1967",
  "description": "Detects the creation of a schtask that executes a file from C:\\Users\\<USER>\\AppData\\Local",
  "level": "high",
  "status": "test",
  "tags": [
    "attack.privilege-escalation",
    "attack.execution",
    "attack.persistence",
    "attack.t1053.005",
    "attack.t1059.001"
  ],
  "category": "process_creation",
  "rule_file": "rules/windows/process_creation/proc_creation_win_schtasks_appdata_local_system.yml",
  "sources": [
    "security_process_creation",
    "sysmon_process_creation"
  ],
  "items": {
    "selection": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "Image",
          "operator": "endswith",
          "values": [
            "\\schtasks.exe"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/Create",
            "/RU",
            "/TR",
            "C:\\Users\\",
            "\\AppData\\Local\\"
          ],
          "match": "all",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "NT AUT",
            " SYSTEM "
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "filter": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "ParentImage",
          "operator": "contains",
          "values": [
            "\\AppData\\Local\\Temp\\",
            "TeamViewer_.exe"
          ],
          "match": "all",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "Image",
          "operator": "endswith",
          "values": [
            "\\schtasks.exe"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/TN TVInstallRestore"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    }
  },
  "condition": {
    "type": "all",
    "children": [
      {
        "type": "item_ref",
        "name": "selection"
      },
      {
        "type": "not",
        "child": {
          "type": "item_ref",
          "name": "filter"
        }
      }
    ]
  }
}