# Windows/Forensicator-Share/rules/windows/sigma/process_creation/proc_creation_win_schtasks_powershell_persistence.json
{
  "title": "Potential Persistence Via Powershell Search Order Hijacking - Task",
  "rule_id": "b66474aa-bd92-4333-a16c-298155b120df",
  "description": "Detects suspicious powershell execution via a schedule task where the command ends with an suspicious flags to hide the powershell instance instead of executeing scripts or commands. This could be a sign of persistence via PowerShell \"Get-Variable\" technique as seen being used in Colibri Loader",
  "level": "high",
  "status": "test",
  "tags": [
    "attack.privilege-escalation",
    "attack.execution",
    "attack.persistence",
    "attack.t1053.005",
    "attack.t1059.001"
  ],
  "category": "process_creation",
  "rule_file": "rules/windows/process_creation/proc_creation_win_schtasks_powershell_persistence.yml",
  "sources": [
    "security_process_creation",
    "sysmon_process_creation"
  ],
  "items": {
    "selection": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "ParentImage",
          "operator": "exact",
          "values": [
            "C:\\WINDOWS\\System32\\svchost.exe"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "ParentCommandLine",
          "operator": "contains",
          "values": [
            "-k netsvcs",
            "-s Schedule"
          ],
          "match": "all",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "endswith",
          "values": [
            " -windowstyle hidden",
            " -w hidden",
            " -ep bypass",
            " -noni"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    }
  },
  "condition": {
    "type": "item_ref",
    "name": "selection"
  }
}