# Windows/Forensicator-Share/rules/windows/sigma/process_creation/proc_creation_win_schtasks_env_folder.json
{
  "title": "Schedule Task Creation From Env Variable Or Potentially Suspicious Path Via Schtasks.EXE",
  "rule_id": "81325ce1-be01-4250-944f-b4789644556f",
  "description": "Detects Schtask creations that point to a suspicious folder or an environment variable often used by malware",
  "level": "medium",
  "status": "test",
  "tags": [
    "attack.privilege-escalation",
    "attack.persistence",
    "attack.execution",
    "attack.t1053.005"
  ],
  "category": "process_creation",
  "rule_file": "rules/windows/process_creation/proc_creation_win_schtasks_env_folder.yml",
  "sources": [
    "security_process_creation",
    "sysmon_process_creation"
  ],
  "items": {
    "selection_1_create": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "Image",
          "operator": "endswith",
          "values": [
            "\\schtasks.exe"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            " /create "
          ],
          "match": "any",
          "windash": true,
          "ignore_case": true
        }
      ]
    },
    "selection_1_all_folders": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            ":\\Perflogs",
            ":\\Users\\All Users\\",
            ":\\Users\\Default\\",
            ":\\Users\\Public",
            ":\\Windows\\Temp",
            "\\AppData\\Local\\",
            "\\AppData\\Roaming\\",
            "%AppData%",
            "%Public%"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "selection_2_parent": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "ParentCommandLine",
          "operator": "endswith",
          "values": [
            "\\svchost.exe -k netsvcs -p -s Schedule"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "selection_2_some_folders": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            ":\\Perflogs",
            ":\\Windows\\Temp",
            "\\Users\\Public",
            "%Public%"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "filter_optional_other": {
      "type": "any",
      "children": [
        {
          "type": "all",
          "children": [
            {
              "type": "field",
              "field": "ParentCommandLine",
              "operator": "contains",
              "values": [
                "unattended.ini"
              ],
              "match": "any",
              "windash": false,
              "ignore_case": true
            }
          ]
        },
        {
          "type": "all",
          "children": [
            {
              "type": "field",
              "field": "CommandLine",
              "operator": "contains",
              "values": [
                "update_task.xml"
              ],
              "match": "any",
              "windash": false,
              "ignore_case": true
            }
          ]
        }
      ]
    },
    "filter_optional_team_viewer": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/Create /TN TVInstallRestore /TR"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "filter_optional_avira_install": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/Create /Xml ",
            "\\Temp\\.CR.",
            "\\Avira_Security_Installation.xml"
          ],
          "match": "all",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "filter_optional_avira_other": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/Create /F /TN",
            "/Xml ",
            "\\Temp\\",
            "Avira_"
          ],
          "match": "all",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            ".tmp\\UpdateFallbackTask.xml",
            ".tmp\\WatchdogServiceControlManagerTimeout.xml",
            ".tmp\\SystrayAutostart.xml",
            ".tmp\\MaintenanceTask.xml"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "filter_optional_klite_codec": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "\\Temp\\",
            "/Create /TN \"klcp_update\" /XML ",
            "\\klcp_update_task.xml"
          ],
          "match": "all",
          "windash": false,
          "ignore_case": true
        }
      ]
    }
  },
  "condition": {
    "type": "all",
    "children": [
      {
        "type": "any",
        "children": [
          {
            "type": "wildcard_ref",
            "mode": "all",
            "pattern": "selection_1_*"
          },
          {
            "type": "wildcard_ref",
            "mode": "all",
            "pattern": "selection_2_*"
          }
        ]
      },
      {
        "type": "not",
        "child": {
          "type": "wildcard_ref",
          "mode": "any",
          "pattern": "filter_optional_*"
        }
      }
    ]
  }
}