# Windows/Forensicator-Share/rules/windows/sigma/process_creation/proc_creation_win_schtasks_persistence_windows_telemetry.json
{
  "title": "Potential Persistence Via Microsoft Compatibility Appraiser",
  "rule_id": "f548a603-c9f2-4c89-b511-b089f7e94549",
  "description": "Detects manual execution of the \"Microsoft Compatibility Appraiser\" task via schtasks.\nIn order to trigger persistence stored in the \"\\AppCompatFlags\\TelemetryController\" registry key.\n",
  "level": "medium",
  "status": "test",
  "tags": [
    "attack.privilege-escalation",
    "attack.execution",
    "attack.persistence",
    "attack.t1053.005"
  ],
  "category": "process_creation",
  "rule_file": "rules/windows/process_creation/proc_creation_win_schtasks_persistence_windows_telemetry.yml",
  "sources": [
    "sysmon_process_creation"
  ],
  "items": {
    "selection_img": {
      "type": "any",
      "children": [
        {
          "type": "all",
          "children": [
            {
              "type": "field",
              "field": "Image",
              "operator": "endswith",
              "values": [
                "\\schtasks.exe"
              ],
              "match": "any",
              "windash": false,
              "ignore_case": true
            }
          ]
        },
        {
          "type": "all",
          "children": [
            {
              "type": "field",
              "field": "OriginalFileName",
              "operator": "exact",
              "values": [
                "schtasks.exe"
              ],
              "match": "any",
              "windash": false,
              "ignore_case": true
            }
          ]
        }
      ]
    },
    "selection_cli": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "run ",
            "\\Application Experience\\Microsoft Compatibility Appraiser"
          ],
          "match": "all",
          "windash": false,
          "ignore_case": true
        }
      ]
    }
  },
  "condition": {
    "type": "wildcard_ref",
    "mode": "all",
    "pattern": "selection_*"
  }
}