# Windows/Forensicator-Share/rules/windows/sigma/process_creation/proc_creation_win_schtasks_one_time_only_midnight_task.json
{
  "title": "Uncommon One Time Only Scheduled Task At 00:00",
  "rule_id": "970823b7-273b-460a-8afc-3a6811998529",
  "description": "Detects scheduled task creation events that include suspicious actions, and is run once at 00:00",
  "level": "high",
  "status": "test",
  "tags": [
    "attack.execution",
    "attack.persistence",
    "attack.privilege-escalation",
    "attack.t1053.005"
  ],
  "category": "process_creation",
  "rule_file": "rules/windows/process_creation/proc_creation_win_schtasks_one_time_only_midnight_task.yml",
  "sources": [
    "sysmon_process_creation"
  ],
  "items": {
    "selection_img": {
      "type": "any",
      "children": [
        {
          "type": "all",
          "children": [
            {
              "type": "field",
              "field": "Image",
              "operator": "contains",
              "values": [
                "\\schtasks.exe"
              ],
              "match": "any",
              "windash": false,
              "ignore_case": true
            }
          ]
        },
        {
          "type": "all",
          "children": [
            {
              "type": "field",
              "field": "OriginalFileName",
              "operator": "exact",
              "values": [
                "schtasks.exe"
              ],
              "match": "any",
              "windash": false,
              "ignore_case": true
            }
          ]
        }
      ]
    },
    "selection_cli": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "wscript",
            "vbscript",
            "cscript",
            "wmic ",
            "wmic.exe",
            "regsvr32.exe",
            "powershell",
            "\\AppData\\"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "selection_time": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "once",
            "00:00"
          ],
          "match": "all",
          "windash": false,
          "ignore_case": true
        }
      ]
    }
  },
  "condition": {
    "type": "wildcard_ref",
    "mode": "all",
    "pattern": "selection_*"
  }
}