# Windows/Forensicator-Share/rules/windows/sigma/process_creation/proc_creation_win_hktl_powersploit_empire_default_schtasks.json
{
  "title": "HackTool - Default PowerSploit/Empire Scheduled Task Creation",
  "rule_id": "56c217c3-2de2-479b-990f-5c109ba8458f",
  "description": "Detects the creation of a schtask via PowerSploit or Empire Default Configuration.",
  "level": "high",
  "status": "test",
  "tags": [
    "attack.execution",
    "attack.persistence",
    "attack.privilege-escalation",
    "attack.s0111",
    "attack.g0022",
    "attack.g0060",
    "car.2013-08-001",
    "attack.t1053.005",
    "attack.t1059.001"
  ],
  "category": "process_creation",
  "rule_file": "rules/windows/process_creation/proc_creation_win_hktl_powersploit_empire_default_schtasks.yml",
  "sources": [
    "security_process_creation",
    "sysmon_process_creation"
  ],
  "items": {
    "selection": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "ParentImage",
          "operator": "endswith",
          "values": [
            "\\powershell.exe",
            "\\pwsh.exe"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "Image",
          "operator": "endswith",
          "values": [
            "\\schtasks.exe"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/Create",
            "powershell.exe -NonI",
            "/TN Updater /TR"
          ],
          "match": "all",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/SC ONLOGON",
            "/SC DAILY /ST",
            "/SC ONIDLE",
            "/SC HOURLY"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    }
  },
  "condition": {
    "type": "item_ref",
    "name": "selection"
  }
}