# Windows/Forensicator-Share/rules/windows/sigma/process_creation/proc_creation_win_schtasks_reg_loader.json
{
  "title": "Scheduled Task Executing Payload from Registry",
  "rule_id": "86588b36-c6d3-465f-9cee-8f9093e07798",
  "description": "Detects the creation of a schtasks that potentially executes a payload stored in the Windows Registry using PowerShell.",
  "level": "medium",
  "status": "test",
  "tags": [
    "attack.privilege-escalation",
    "attack.execution",
    "attack.persistence",
    "attack.t1053.005",
    "attack.t1059.001"
  ],
  "category": "process_creation",
  "rule_file": "rules/windows/process_creation/proc_creation_win_schtasks_reg_loader.yml",
  "sources": [
    "sysmon_process_creation"
  ],
  "items": {
    "selection_img": {
      "type": "any",
      "children": [
        {
          "type": "all",
          "children": [
            {
              "type": "field",
              "field": "Image",
              "operator": "endswith",
              "values": [
                "\\schtasks.exe"
              ],
              "match": "any",
              "windash": false,
              "ignore_case": true
            }
          ]
        },
        {
          "type": "all",
          "children": [
            {
              "type": "field",
              "field": "OriginalFileName",
              "operator": "exact",
              "values": [
                "schtasks.exe"
              ],
              "match": "any",
              "windash": false,
              "ignore_case": true
            }
          ]
        }
      ]
    },
    "selection_cli_create": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/Create"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "selection_cli_get": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "Get-ItemProperty",
            " gp "
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "selection_cli_hive": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "HKCU:",
            "HKLM:",
            "registry::",
            "HKEY_"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "filter_main_encoding": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "FromBase64String",
            "encodedcommand"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    }
  },
  "condition": {
    "type": "all",
    "children": [
      {
        "type": "wildcard_ref",
        "mode": "all",
        "pattern": "selection_*"
      },
      {
        "type": "not",
        "child": {
          "type": "wildcard_ref",
          "mode": "any",
          "pattern": "filter_*"
        }
      }
    ]
  }
}