# Windows/Forensicator-Share/rules/windows/sigma/process_creation/proc_creation_win_schtasks_guid_task_name.json
{
  "title": "Suspicious Scheduled Task Name As GUID",
  "rule_id": "ff2fff64-4cd6-4a2b-ba7d-e28a30bbe66b",
  "description": "Detects creation of a scheduled task with a GUID like name",
  "level": "medium",
  "status": "test",
  "tags": [
    "attack.privilege-escalation",
    "attack.persistence",
    "attack.execution",
    "attack.t1053.005"
  ],
  "category": "process_creation",
  "rule_file": "rules/windows/process_creation/proc_creation_win_schtasks_guid_task_name.yml",
  "sources": [
    "security_process_creation",
    "sysmon_process_creation"
  ],
  "items": {
    "selection_img": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "Image",
          "operator": "endswith",
          "values": [
            "\\schtasks.exe"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        },
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/Create "
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "selection_tn": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "/TN \"{",
            "/TN '{",
            "/TN {"
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    },
    "selection_end": {
      "type": "all",
      "children": [
        {
          "type": "field",
          "field": "CommandLine",
          "operator": "contains",
          "values": [
            "}\"",
            "}'",
            "} "
          ],
          "match": "any",
          "windash": false,
          "ignore_case": true
        }
      ]
    }
  },
  "condition": {
    "type": "wildcard_ref",
    "mode": "all",
    "pattern": "selection_*"
  }
}